# openwrt/openmmg.init
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=1
STOP=99

start_service() {
	procd_open_instance
	procd_set_param command /usr/bin/openmmg -c /etc/config/openmmg
	procd_set_param respawn ${respawn_threshold:-3600} ${respawn_timeout:-5} ${respawn_retry:-0}
	procd_close_instance
}

service_triggers() {
	procd_add_reload_trigger "openmmg"
}
